# docker/service/cron/run
#!/bin/sh

exec crond -f -d 0 -l 2 -L /var/log/cron.log -c /var/spool/cron/crontabs